# nasa/nasa_restart.service
[Unit]
Description=Restart Nasa

[Service]
Type=oneshot
ExecStart=systemctl restart nasa.service
