# nasa/nasa_restart.timer
[Unit]
Description=Restart nasa timer

[Timer]
OnCalendar=*-*-* 17:30:00
Persistent=true

[Install]
WantedBy=timers.target
